# server/packaging/weatherstation.service
[Service]
Restart=always
ExceStart=/opt/weatherstation/server
WorkingDirectory=/opt/weatherstation
User=pi

