# server/packaging/weatherstation.service
[Service]
Restart=always
ExecStart=/opt/weatherstation/server
WorkingDirectory=/opt/weatherstation
User=pi

[Install]
WantedBy=multi-user.target

